# tea-server-platform/agent/agent.service
[Unit]
Description=Tea Server Platform Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/tea-agent/agent.py
# Create /etc/tea-agent/agent.env before starting this service:
# AGENT_API_KEY=<platform agent_api_key>
# VIRT_TYPE=lxd
# PLATFORM_URL=https://your-platform-domain.com
EnvironmentFile=/etc/tea-agent/agent.env
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
